# etc/service.d/agetty-generic/run
#!/bin/sh

[ -x /sbin/getty -o -x /bin/getty ]   && exec getty $@
[ -x /sbin/agetty -o -x /bin/agetty ] && exec agetty $@